# container_builds/Podman/generic/images/archlinux_base_with_xserver/data/misc-init.service
[Unit]
Description=Random init stuff

[Service]
ExecStart=/misc-init.sh

[Install]
WantedBy=multi-user.target